# supervisord/barcode-zmq-sub.conf
[program:barcode_zmq_sub]
command=python3 -u barcode_zmq_sub.py
directory=/home/pi/raspberrypi4/websocket
stdout_logfile=/home/pi/raspberrypi4/websocket/barcode_zmq_sub.log
redirect_stderr=true
